# ansible/roles/cameraupload/files/prune_cam_archive.timer
[Unit]
Description=Run camera archive prune every night

[Timer]
Unit=prune_cam_archive.service
OnCalendar=01:00

[Install]
WantedBy=timers.target
