# ansible/roles/cameraupload/files/prune_cam_archive.service
[Unit]
Description=Prune camera archive

[Service]
Type=oneshot
ExecStart=/usr/local/lib/prune_cam_archive.sh
