# ansible/roles/minio/files/minio_prune_cam_archive.service
[Unit]
Description=Prune minio camera archive

[Service]
Type=oneshot
ExecStart=/usr/local/lib/minio_prune_cam_archive.sh
