# ansible/roles/dokuwiki/files/dokuwiki_prune_logs.service
[Unit]
Description=Prune dokuwiki logs

[Service]
Type=oneshot
ExecStart=/usr/local/lib/dokuwiki_prune_logs.sh
