# ansible/roles/backup/files/backup.service
[Unit]
Description=Backup services

[Service]
Type=oneshot
ExecStart=/usr/local/lib/backup
