# ansible/roles/backup/files/backup.timer
[Unit]
Description=Run backup every night

[Timer]
Unit=backup.service
OnCalendar=04:00

[Install]
WantedBy=timers.target
