# ansible/roles/acsgw/files/acs_prune_logs.service
[Unit]
Description=Prune ACS logs

[Service]
Type=oneshot
ExecStart=/usr/local/lib/acs_prune_logs.sh
