# ansible/roles/acsgw/files/acs_prune_logs.timer
[Unit]
Description=Run ACS prune logs every night

[Timer]
Unit=acs_prune_logs.service
OnCalendar=02:00

[Install]
WantedBy=timers.target
